# proprietary/vendor/ramdisk/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:-/bin/sh
tty1::respawn:-/bin/sh
